# cron/crontab
0 22 * * * /opt/renew_certs.sh >> /var/log/cron.log 2>&1